# alpine/rootfs/etc/supervisor.d/wdmrc.ini
[program:wdmrc]
user=runner
command=mono /opt/home/wdmrc/wdmrc.exe -p 8081 -h http://* 
autostart=true
startsecs=10
autorestart=true
exitcodes=0,2
stopsignal=TERM
stdout_logfile=NONE
redirect_stderr=true
